# docker-compose/crontab
* * * * * cd /var/www && php artisan schedule:run >> /dev/stdout 2>&1
