# docker-compose/crontab
* * * * * cd /var/www/html && php artisan schedule:run >> /dev/stdout 2>&1
